# project/config/production/supervisor/swoole_api_frame.conf
[program:swoole_swoole_api_frame]
command=/usr/bin/php /var/www/swoole_api_frame/public/index.php
stopwaitsecs=5
restart=always
user=root
