# project/config/development/supervisor/swoole_developer.conf
[program:swoole_developer]
command=/bin/bash /var/www/swoole_api_frame/project/tool/dev_server_swoole_watch.sh
stopwaitsecs=5
restart=always
user=root
